# image/configs/supervisor.d/snmpd.ini
[program:snmpd]
command =/bin/bash -c "exec /usr/sbin/snmpd -f /etc/snmp/snmpd.conf -p /var/run/snmpd.pid -I -smux 0.0.0.0"
process_name=%(program_name)s
log_stdout=true
log_stderr=true
directory=/
numprocs=1
priority=958
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
user=root